# service/rabbitmq/run
#!/bin/sh

mkdir -p /srv/rabbitmq/log
chown -R rabbitmq:rabbitmq /srv/rabbitmq/log
exec rabbitmq-server

